# services/moltbook-vote-refresh.service
[Unit]
Description=Moltbook full vote refresh (sort=top, all pages)
After=postgresql.service
Requires=postgresql.service
# Stop the scraper daemon to avoid double rate-limit usage.
# Scraper saves cursor state, so it resumes cleanly after restart.
Conflicts=moltbook-scraper.service moltbook-sync.service

[Service]
Type=oneshot
User=%u
WorkingDirectory=/opt/moltbook
ExecStart=/opt/moltbook/.venv/bin/python src/scrape.py --refresh-votes-bulk top --refresh-votes-pages 25000 --skip-posts --skip-comments --skip-submolts
ExecStartPost=+/bin/systemctl start moltbook-scraper.service
ExecStartPost=+/bin/systemctl start moltbook-sync.service
TimeoutStartSec=infinity

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moltbook-vote-refresh
